# glue-bootstrap.service.in
[Service]
ExecStart=@sbindir@/glue-bootstrap
